# .ebextensions/scripts/conf/supervisord.conf
[inet_http_server]
port = 0.0.0.0:49153
username = username
password = password

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:nrgi_queue]
command=php artisan queue:work --env=production --daemon --sleep=3 --tries=3
user=webapp
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor.log
redirect_stderr=true

[program:elasticsearch_queue]
command=php artisan queue:work --queue=elastic_search --env=production --daemon --sleep=3 --tries=3
user=webapp
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor_elasticsearch.log
redirect_stderr=true

[program:queue-mail]
command=php artisan queue:work --queue=queue-mail --env=production --daemon --sleep=3 --tries=3
user=webapp
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor_email.log
redirect_stderr=true

[program:queue-mturk]
command=php artisan queue:work --queue=mturk --env=production --daemon --sleep=3 --tries=3
user=webapp
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor_email.log
redirect_stderr=true

[program:queue-contract_bulk_create]
command=php artisan queue:work --queue=contract_bulk_create --env=production --daemon --sleep=3 --tries=3
user=webapp
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor_email.log
redirect_stderr=true

[program:queue-contract_download]
command=php artisan queue:work --queue=contract_download --env=production --daemon --sleep=3 --tries=3
user=webapp
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor_email.log
redirect_stderr=true

[program:beanstalkd]
command=/usr/bin/beanstalkd
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor_beanstalk.log
redirect_stderr=true

