# conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:remote-syslog]
command = /usr/local/bin/remote_syslog -D
user = root
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/remote-syslog.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:nrgi_queue]
command=php artisan queue:listen --env=production --timeout=0 --tries=2
user=root
directory=/var/www/rc-admin
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/supervisor_nrgi.log
redirect_stderr=true

[program:elasticsearch_queue]
command=php artisan queue:listen --queue=elastic_search --env=production --timeout=0 --tries=2
user=root
directory=/var/www/rc-admin
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/supervisor_elasticsearch.log
redirect_stderr=true

[program:queue-mail]
command=php artisan queue:listen --queue=queue-mail --env=production --timeout=0 --tries=2
user=root
directory=/var/www/rc-admin
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/supervisor_email.log
redirect_stderr=true

[program:queue-mturk]
command=php artisan queue:work --queue=mturk --env=production --daemon --sleep=3 --tries=3
user=root
directory=/var/www/rc-admin
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/supervisor_email.log
redirect_stderr=true

[program:queue-contract_bulk_create]
command=php artisan queue:work --queue=contract_bulk_create --env=production --daemon --sleep=3 --tries=3
user=root
directory=/var/www/rc-admin
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/supervisor_email.log
redirect_stderr=true

[program:queue-contract_download]
command=php artisan queue:work --queue=contract_download --env=production --daemon --sleep=3 --tries=3
user=root
directory=/var/www/rc-admin
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/supervisor_email.log
redirect_stderr=true

